# modules/portal-metrics/src/main/resources/apiportal.cron
*/15 * * * * /opt/SecureSpan/ApiPortal/portal.sh
